# modules/tftp/templates/tftp.service.erb
[Unit]
Description=Tftp Server
Requires=tftp.socket
Documentation=man:in.tftpd

[Service]
ExecStart=/usr/sbin/in.tftpd --timeout 300 -s /tftpboot
StandardInput=socket

[Service]
Restart=always
RestartSec=5

[Install]
Also=tftp.socket
